# assistant/timers/templates/task.timer.tmpl
[Unit]
Description=Timer for Claude Assistant Task: {{NAME}}

[Timer]
OnCalendar={{SCHEDULE}}
Persistent=true

[Install]
WantedBy=timers.target
